# deploy/aprilvision-dashboard.service
[Unit]
Description=AprilVision 3.2 - Dashboard Service
Documentation=https://github.com/TomAs-1226/AprilVIsion-2.0
After=photonvision.service
Requires=photonvision.service

[Service]
Type=simple
User=photonvision
Group=photonvision
WorkingDirectory=/opt/aprilvision
# NOTE: This file is a template. setup.sh --deploy generates the actual service
# with the correct --team flag. Edit /etc/systemd/system/aprilvision-dashboard.service
# for the live version, or re-run: sudo ./setup.sh --deploy --team YOUR_TEAM
ExecStart=/usr/bin/python3 /opt/aprilvision/bridge.py --port 5801 --engine-port 5800 --web-dir /opt/aprilvision/web
Restart=always
RestartSec=3
TimeoutStartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aprilvision-dashboard

[Install]
WantedBy=multi-user.target
