# deploy/frc_vision.service
[Unit]
Description=FRC AprilTag Vision Coprocessor
Documentation=https://github.com/team/frc-vision

# Start after network target but DO NOT block boot waiting for network
# Service handles network unavailability gracefully with internal retries
After=network.target
Wants=network-online.target

# Don't wait for network-online - we handle it in-app
DefaultDependencies=no
After=sysinit.target

[Service]
Type=simple

# Run as dedicated user (created by install script)
User=frcvision
Group=frcvision

# Working directory
WorkingDirectory=/opt/frc-vision

# Main executable with fast-start mode flag
ExecStart=/opt/frc-vision/frc_vision /opt/frc-vision/config/config.yml --fast-start

# CRITICAL: Fast restart for reliability
Restart=always
RestartSec=0.5

# Startup timeout - fail fast if something is really broken
TimeoutStartSec=30
TimeoutStopSec=10

# Nice priority (slight boost, but safe)
Nice=-5

# CPU scheduling - use batch for better throughput without blocking system
# CPUSchedulingPolicy=batch is safer than fifo for boot
CPUSchedulingPolicy=other
CPUSchedulingPriority=0

# Resource limits
LimitMEMLOCK=infinity
LimitRTPRIO=99
LimitNOFILE=65536

# Memory limits to prevent runaway
MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=frc-vision

# Environment - minimal for fast start
Environment="HOME=/opt/frc-vision"
Environment="DISPLAY="
Environment="LC_ALL=C"

# Don't send SIGHUP on terminal disconnect
IgnoreSIGPIPE=yes

# Watchdog - restart if frozen (optional, comment out if issues)
# WatchdogSec=30

# systemd-notify integration (optional)
# Type=notify
# NotifyAccess=main

[Install]
WantedBy=multi-user.target
